# bashconf/startupconf/startup.service

[Unit]
Description=Startup file

[Service]
ExecStart=/bin/bash /usr/local/sbin/startup.sh

[Install]
WantedBy=multi-user.target
